# debian/cuckoo-restapi.service
[Service]
User=cuckoo
EnvironmentFile=-/etc/default/cuckoo
ExecStart=/usr/bin/python /usr/share/cuckoo/utils/api.py $CUCKOO_RESTAPI_ARGS
