# debian/cuckoo-distributed.service
[Service]
User=cuckoo
EnvironmentFile=-/etc/default/cuckoo
ExecStart=/usr/bin/python /usr/share/cuckoo/utils/dist.py $CUCKOO_DISTRIBUTED_ARGS
